# autostart/ai-usage-tray.desktop
[Desktop Entry]
Type=Application
Name=AI Usage Tray
Comment=Tracks AI coding assistant usage and displays remaining limits
Exec=ai-usage-tray
Icon=utilities-system-monitor
Terminal=false
Categories=Utility;System;
StartupNotify=false
X-GNOME-Autostart-enabled=true
